# iiii.conf
[program:iiii]
command=python3 /home/curator2/iiii/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/iiii/test.err.log
stdout_logfile=/home/curator2/iiii/test.out.log